# docker/local/etc/services.d/laravel-queue/run
#!/usr/bin/with-contenv sh
# Make sure all queue will be run as `nginx` user.

s6-setuidgid nginx ./.default-queue.sh